# superhost-api.service
[Unit]
Description=Superhost API Service
After=network.target postgresql.service

[Service]
Type=simple
User=jonathan
WorkingDirectory=/home/jonathan/superhost/api
ExecStart=/usr/bin/node --require /home/jonathan/superhost/api/node_modules/tsx/dist/preflight.cjs --import file:///home/jonathan/superhost/api/node_modules/tsx/dist/loader.mjs src/index.ts
Restart=always
Environment=NODE_ENV=production
# Inject .env values here or use EnvironmentFile
EnvironmentFile=/home/jonathan/superhost/api/.env

[Install]
WantedBy=multi-user.target
